# my-script/debian/my-script/etc/systemd/system/myscrit.service
[Unit]
Description=MyScript

[Service]
ExecStart= /bin/bash /etc/systemd/system/final.sh

[Install]
WantedBy=mytimer.target
