# package_build_dir/my_package_name/etc/systemd/system/myscrit.service
[Unit]
Description=MyScript

[Service]
ExecStart= /bin/bash ~/final.sh

[Install]
WantedBy=mytimer.target
